# web-portal/backend/scripts/systemd/pooldash_portal.service
[Unit]
Description=PoolDash Portal Backend
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=pooldash
Group=pooldash
WorkingDirectory=/opt/pooldash/portal/backend
EnvironmentFile=/opt/pooldash/portal/backend/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pooldash/portal/backend/uploads

[Install]
WantedBy=multi-user.target
