# pi-software/PoolDash_v6/scripts/settings_backup.timer
[Unit]
Description=Backup PoolAIssistant settings daily

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
