# pi-software/PoolDash_v6/scripts/watchdog.service
[Unit]
Description=PoolAIssistant Service Watchdog
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/python3 /opt/PoolAIssistant/app/scripts/watchdog.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
