# web-portal/backend/scripts/cron/cleanup_readings.cron
0 3 * * * cd /opt/pooldash/portal/backend && /usr/bin/npm run cleanup-readings >> /var/log/pooldash_cleanup.log 2>&1
