# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_weekly_reboot.timer
[Unit]
Description=Weekly Reboot for PoolAIssistant (clears memory fragmentation)

[Timer]
OnCalendar=Sun 04:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
