# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_scheduled_reboot.service
[Unit]
Description=PoolAIssistant Scheduled Daily Reboot
Documentation=man:systemd-reboot.service(8)

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now "Scheduled daily reboot for system health"
