# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_weekly_reboot.service
[Unit]
Description=Weekly Reboot for PoolAIssistant
Documentation=man:systemd-reboot.service(8)

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now "Weekly scheduled reboot for system health"
