# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_device_sync.timer
[Unit]
Description=PoolAIssistant Device Sync Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
AccuracySec=1min
Persistent=true
Unit=poolaissistant_device_sync.service

[Install]
WantedBy=timers.target
