# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_scheduled_reboot.timer
[Unit]
Description=PoolAIssistant Daily Scheduled Reboot

[Timer]
# Default time - will be overwritten by configure_scheduled_reboot.sh
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
