# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_device_sync.service
[Unit]
Description=PoolAIssistant Device Provisioning and Upload
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/PoolAIssistant/app
EnvironmentFile=/etc/PoolAIssistant/poolaissistant.env
ExecStart=/opt/PoolAIssistant/venv/bin/python /opt/PoolAIssistant/app/scripts/device_upload.py

[Install]
WantedBy=multi-user.target
