# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_provision.service
[Unit]
Description=PoolAIssistant Auto-Provisioning
After=network-online.target
Wants=network-online.target
Before=poolaissistant_ui.service
ConditionPathExists=!/opt/PoolAIssistant/data/.provisioned

[Service]
Type=oneshot
User=poolai
Group=poolai
WorkingDirectory=/opt/PoolAIssistant/app
ExecStart=/usr/bin/python3 /opt/PoolAIssistant/app/scripts/auto_provision.py
ExecStartPost=/bin/touch /opt/PoolAIssistant/data/.provisioned
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Retry on failure
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
