# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_data_cleanup.timer
[Unit]
Description=PoolAIssistant Data Cleanup Timer

[Timer]
# Run daily at 3 AM
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
