# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_ports.service
[Unit]
Description=PoolAIssistant Port Configuration
After=network.target
Before=poolaissistant_ui.service

[Service]
Type=oneshot
ExecStart=/opt/PoolAIssistant/app/scripts/ensure_ports.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
