# pi-software/PoolDash_v6/scripts/systemd/poolaissistant_button.service
[Unit]
Description=PoolAIssistant Physical Button Handler
After=multi-user.target
Wants=network.target

[Service]
Type=simple
ExecStart=/opt/PoolAIssistant/venv/bin/python /opt/PoolAIssistant/app/scripts/button_handler.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
